# docker/agiper-agiper/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld_safe --binlog-format=MIXED

[program:teamcity]
command=/opt/TeamCity/bin/runAll.sh start